# examples/external_display/autostart/external-display-acm0.desktop
[Desktop Entry]
Type=Application
Name=External Display ACM0
Exec=/usr/bin/bash -c "$HOME/external-display-acm0.sh"
